# systemd/system/caches-manager.service
[Unit]
Description=Caches manager: create system caches

[Service]
Type=oneshot
ExecStart=caches-manager --system

[Install]
WantedBy=basic.target
