# systemd/user/caches-manager.service
[Unit]
Description=Caches manager: create user caches

[Service]
Type=oneshot
ExecStart=caches-manager --user

[Install]
WantedBy=basic.target
